# scripts/kiosk.service
[Unit]
Description=Skylight-style kiosk app
After=network-online.target
Wants=network-online.target

[Service]
User=pi
Environment=DISPLAY=:0
WorkingDirectory=/home/pi/hub
ExecStart=/home/pi/hub/scripts/start-kiosk.sh
Restart=always

[Install]
WantedBy=graphical.target

