# scripts/kiosk.service
[Unit]
Description=Skylight-style kiosk app
After=network-online.target
Wants=network-online.target

[Service]
User=pi
Environment=DISPLAY=:0
WorkingDirectory=/home/pi/hub
ExecStart=/usr/bin/chromium --noerrdialogs --disable-infobars --kiosk --app=http://localhost --touch-events --enable-pinch --enable-viewport --enable-threaded-scrolling --enable-smooth-scrolling --overscroll-history-navigation=0 --check-for-update-interval=31536000 --user-agent="Mozilla/5.0 (Linux; Android 13; Tablet) AppleWebKit/537.36 (KHTML, like Gecko) Chrome/120.0.0.0 Safari/537.36"
Restart=always

[Install]
WantedBy=graphical.target

